# .config/systemd/user/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=xremap /home/rnakamura/.config/xremap.yml --watch
